# debian/fastcgi-daemon2-init.files
/usr/bin/fastcgistart2.sh
/etc/init.d/fastcgi-daemon2

